# files/base-server/etc/systemd/user/podman-auto-update.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=
OnCalendar=Sun *-*-* 04:30:00
Persistent=true
EOF
