# files/base/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]

OnBootSec=
OnUnitInactiveSec=
OnCalendar=


OnCalendar=Sun *-*-* 04:00:00


Persistent=true